# debian/flatpak-automatic.timer
[Unit]
Description=Run daily Flatpak update check

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
